# docker/root/etc/services.d/tdarr_node/run
#!/usr/bin/with-contenv bash

echo "Starting Tdarr_Node"

cd /app/Tdarr_Node || exit 1

exec s6-setuidgid Tdarr ./Tdarr_Node
