# docker/tdarr/root/etc/services.d/mongodb/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc mongod
