# docker/tdarr/root/etc/services.d/tdarr/run
#!/usr/bin/with-contenv bash

exec \
      s6-setuidgid abc node --max-old-space-size=16384 /app/bundle/main.js 
